# supervisord.conf
[supervisord]
nodaemon = true
 
[program:scrapyd]
command = scrapyd
autostart = true
autorestart = true
 
[program:spiderkeeper]
command = spiderkeeper
autostart = true
autorestart = true
